# files/gunicorn3@.socket
[Unit]
Description=gunicorn socket for %i

[Socket]
ListenStream=/run/gunicorn/%i.socket

[Install]
WantedBy=sockets.target