# files/gunicorn3@.socket
[Unit]
Description=gunicorn socket for ${appname}

[Socket]
ListenStream=/run/gunicorn/${appname}.socket

[Install]
WantedBy=sockets.target